# config/systemd/user/mojibar.service
[Unit]
Description=Start Mojibar as a service

[Service]
ExecStart=bash -c "mojibar"

[Install]
WantedBy=default.target
